# sign_user/sign.service
[Unit]
Description=ELF Sign Service

[Service]
Type=simple
ExecStart=/usr/bin/user /etc/public.pem

[Install]
WantedBy=multi-user.target